# antispam-daemon.conf
# antispam-daemon - sa-learn daemon

description	"Monitors Spam from Dovecot-Spam"

start on runlevel [2345]

respawn
respawn limit 10 5

exec antispam-daemon
